# supervisor/chat.conf
[program:mattermost]
directory=/home/web/repo/project/chat.pkimber.net/live
command=/home/web/repo/project/chat.pkimber.net/live/bin/platform
startsecs=5
autorestart=true
site={{ site }}
