# supervisor/chat.conf
[program:mattermost]
autorestart=true
directory=/home/web/repo/project/{{ domain }}/live
command=/home/web/repo/project/{{ domain }}/live/bin/platform
startsecs=5
user=web
