# supervisor/chat.conf
[program:mattermost]
directory=/home/web/repo/project/{{ site }}/live
command=/home/web/repo/project/{{ site }}/live/bin/platform
startsecs=5
autorestart=true
