# supervisor/uwsgi.conf
[program:uwsgi]
user=web
command=uwsgi --emperor /home/web/repo/uwsgi/vassals
autostart=true
autorestart=true
