# supervisor/uwsgi.conf
[program:uwsgi]
user=web
command=sh /home/web/opt/runinenv.sh /home/web/repo/uwsgi/venv_uwsgi uwsgi --emperor /home/web/repo/uwsgi/vassals/
autostart=true
autorestart=true
