# supervisor/uwsgi.conf
[program:uwsgi]
user=web
command=/opt/uwsgi/uwsgi --emperor /home/web/repo/uwsgi/vassals --binary-path /opt/uwsgi/uwsgi
autostart=true
autorestart=true
